# linnovate.service.j2
# linnovate.service
[Unit]
Description = Linnovate Initialization Service
After       = syslog.target

[Service]
Environment = SOFTWARE={{ software_name }}
ExecStart   = /opt/linnovate/post.sh $SOFTWARE

[Install]
WantedBy=default.target
